# mac/com.focusboard.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.focusboard.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/josiahlachance/Documents/Projects/Claude/terminal/focusboard/mac/sync_runner.py</string>
    </array>

    <key>StartInterval</key>
    <integer>120</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/josiahlachance/.claude/pi/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/josiahlachance/.claude/pi/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin</string>
        <key>HOME</key>
        <string>/Users/josiahlachance</string>
    </dict>
</dict>
</plist>
